# aerobake/aeroscale/_etc_init.d_amazon-cloudwatch-agent
#!/sbin/openrc-run

description="Amazon CloudWatch Agent"
daemon_bin="/opt/aws/amazon-cloudwatch-agent/bin/amazon-cloudwatch-agent"
# Ensure you have generated this TOML once or created it manually
config_file="/opt/aws/amazon-cloudwatch-agent/etc/amazon-cloudwatch-agent.toml"

depend() {
    need net
    after firewall
}

start_pre() {
    # 1. Get IMDSv2 Token
    TOKEN=$(curl -s -X PUT "http://169.254.169.254/latest/api/token" \
        -H "X-aws-ec2-metadata-token-ttl-seconds: 21600")

    if [ -z "$TOKEN" ]; then
        eerror "Failed to fetch IMDSv2 Token"
        return 1
    fi

    # 2. Get Instance ID using the Token
    INSTANCE_ID=$(curl -s -H "X-aws-ec2-metadata-token: $TOKEN" \
        "http://169.254.169.254/latest/meta-data/instance-id")

    if [ -z "$INSTANCE_ID" ]; then
        eerror "Failed to fetch Instance ID"
        return 1
    fi

    # Export for the start function
    export INSTANCE_ID
}

start() {
    ebegin "Starting Amazon CloudWatch Agent for $INSTANCE_ID"
    # --env passes the variable so the agent can resolve $INSTANCE_ID in the TOML
    start-stop-daemon --start --exec "$daemon_bin" \
        --env INSTANCE_ID="$INSTANCE_ID" \
        --pidfile /var/run/amazon-cloudwatch-agent.pid \
        --make-pidfile --background \
        --user root \
        -- -config "$config_file"
    eend $?
}

stop() {
    ebegin "Stopping Amazon CloudWatch Agent"
    start-stop-daemon --stop --pidfile /var/run/amazon-cloudwatch-agent.pid
    eend $?
}

